# supervisord.conf
[supervisord]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:webserver]
command = gunicorn --workers=3 --max-requests 10000 --max-requests-jitter 10 --bind 0.0.0.0:8000 watchtower.wsgi
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:celery_worker_beat]
command = celery -A watchtower worker --beat -l INFO -c 1 -Ofair -Q beat
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:celery_latest_blockheight_getter]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q latest_blockheight_getter --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:celery_deposit_filter]
command = celery -A watchtower worker -l INFO -c 8 -Ofair -Q deposit_filter --autoscale=1,8 --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:celery_save_record]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q save_record --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_slpbitcoin]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q slpbitcoinsocket
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:slpsocket_fountainhead]
command = python manage.py slpsocket_fountainhead
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:slpstream_fountainhead]
command = python manage.py slpstream_fountainhead
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_bitsocket]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q bitsocket
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true



[program:celery_bitdbquery]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q bitdbquery --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:celery_bitcoincash_tracker]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q bitcoincash_tracker --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_first_blockheight_scanner]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q first_blockheight_scanner --max-tasks-per-child=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_client_acknowledgement]
command = celery -A watchtower worker -l INFO -c 8 -Ofair -Q client_acknowledgement --autoscale=1,8 --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true



[program:celery_openfromredis]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q openfromredis --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_suspendtoredis]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q suspendtoredis --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program: celery_slpdb_token_scanner]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q slpdb_token_scanner --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program: celery_second_blockheight_scanner]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q second_blockheight_scanner --max-tasks-per-child=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program: celery_checktransaction]
command = celery -A watchtower worker -l INFO -c 8 -Ofair -Q checktransaction --autoscale=1,8 --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program: celery_send_telegram_message]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q send_telegram_message --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program: celery_send_slack_message]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q send_slack_message --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program: celery_spicebot_subscription]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q spicebot_subscription
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program: celery_updates]
command = celery -A watchtower worker -l INFO -c 1 -Ofair -Q updates
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true
