# supervisord.conf
[supervisord]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:http_server]
command=gunicorn watchtower.wsgi:application --workers 8 --threads 8 --timeout 15 --max-requests 10000 --bind 0.0.0.0:8000
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:websocket_server]
command=uvicorn --host 0.0.0.0 --port 9000 --loop uvloop --http h11 --interface asgi3 --no-access-log watchtower.asgi:application
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_client_acknowledgement]
command=celery -A watchtower worker -n worker1 -l INFO -Ofair -Q client_acknowledgement --autoscale=1,8 --max-tasks-per-child=100
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:bchn_zmq_listener]
directory = /code
command=python manage.py bchn_zmq_listener
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:stablehedge_tx_queue]
command=python manage.py redemption_contract_txs_queue
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:celery_worker_beat]
command=celery -A watchtower beat
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_save_record]
command=celery -A watchtower worker -n worker4 -l INFO -Ofair -Q save_record --max-tasks-per-child=100  --autoscale=1,8
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_post_save_record]
command=celery -A watchtower worker -n worker5 -l INFO -Ofair -Q post_save_record --max-tasks-per-child=100  --autoscale=1,10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_contract_history]
command=celery -A watchtower worker -n contract_history.worker -l INFO -Ofair -Q contract_history --max-tasks-per-child=100  --autoscale=1,5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_wallet_history_1]
command=celery -A watchtower worker -n worker6 -l INFO -Ofair -Q wallet_history_1 --max-tasks-per-child=100  --autoscale=1,5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_wallet_history_2]
command=celery -A watchtower worker -n worker7 -l INFO -Ofair -Q wallet_history_2 --max-tasks-per-child=100  --autoscale=1,5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_get_latest_block]
command=celery -A watchtower worker -n worker8 -n worker1  -l INFO -c 1 -Ofair -Q get_latest_block --max-tasks-per-child=100
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_manage_blocks]
command=celery -A watchtower worker -n worker9 -l INFO -c 1 -Ofair -Q manage_blocks --max-tasks-per-child=100
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_query_transaction]
command= celery -A watchtower worker -n worker10 -l INFO -Ofair -Q query_transaction --max-tasks-per-child=100 --autoscale=1,8
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:celery_get_utxos]
command=celery -A watchtower worker -n worker11 -l INFO -Ofair -Q get_utxos --max-tasks-per-child=100 --autoscale=1,4
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_token_metadata]
command=celery -A watchtower worker -n worker12 -l INFO -Ofair -Q token_metadata --max-tasks-per-child=100
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_txn_broadcast]
command=celery -A watchtower worker -n worker13 -l INFO -Ofair -Q broadcast --max-tasks-per-child=100 --autoscale=1,10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:celery_rescan_utxos]
command=celery -A watchtower worker -n worker15 -l INFO -Ofair -Q rescan_utxos --max-tasks-per-child=100 --autoscale=1,5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:anyhedge__tasks_queue1]
command = celery -A watchtower worker -n worker20 -l INFO -Ofair -Q anyhedge__funding_parser
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:anyhedge__tasks_queue2]
command = celery -A watchtower worker -n worker21 -l INFO -Ofair -Q anyhedge__price_oracle,anyhedge__settlement_updates
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:stablehedge_tasks]
command=celery -A watchtower worker -n worker27 -l INFO -Ofair -Q stablehedge__treasury_contract --max-tasks-per-child=100 --autoscale=1,10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


; [program:chat__mqtt_listener]
; command = python manage.py mqtt_listener
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
; stderr_logfile_maxbytes=0
; stopasgroup=true


[program:celery_shift_expiration]
command = celery -A watchtower worker -n worker22 -l INFO -Ofair -Q ramp__shift_expiration
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:paytacagifts__monitor]
command = celery -A watchtower worker -n worker24 -l INFO -Ofair -Q monitor-gifts --max-tasks-per-child=100 --autoscale=1,2
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:populate_token_addresses]
command=celery -A watchtower worker -n worker25 -l INFO -Ofair -Q populate_token_addresses --max-tasks-per-child=100  --autoscale=1,2
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:mempool_processing_throttled]
command=celery -A watchtower worker -n worker26 -l INFO -Ofair -Q mempool_processing_throttled --max-tasks-per-child=1000  --autoscale=1,2
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:mempool_processing_fast]
command=celery -A watchtower worker -n worker26 -l INFO -Ofair -Q mempool_processing_fast --max-tasks-per-child=1000  --autoscale=1,20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


[program:mempool_listener]
command=python manage.py mempool_listener
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:js_express_server]
command=node /code/main/js/server.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:anyhedge_js_scripts_server]
command = node /code/anyhedge/js/src/server.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:cts_express_server]
command=npm start -C /code/cts/js
environment=CTS_APP_ENV=%(ENV_CTS_APP_ENV)s
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:stablehedge_js_scripts_server]
command = node /code/stablehedge/js/src/server.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

[program:rampp2p__marketprices]
command = celery -A watchtower worker -n rampp2p__marketprices -l INFO -Ofair -Q rampp2p__marketprices
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:rampp2p__contract_execution]
command = celery -A watchtower worker -n rampp2p__contract_execution -l INFO -Ofair -Q rampp2p__contract_execution
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:rampp2p__cancel_expired_orders]
command = celery -A watchtower worker -n rampp2p__cancel_expired_orders -l INFO -Ofair -Q rampp2p__cancel_expired_orders
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:paytacapos__cashout]
command = celery -A watchtower worker -n paytacapos__cashout -l INFO -Ofair -Q paytacapos__cashout
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true

[program:multisig_express_server]
command=node /code/multisig/js/server.js
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true


# [program:chat__tasks_queue1]
# command = celery -A watchtower worker -n worker23 -l INFO -Ofair -Q chat__notifications --max-tasks-per-child=100 --autoscale=1,2
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true


# [program:chat__mqtt_listener]
# command = python manage.py mqtt_listener
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true

# [program:sbch_blockheader_stream]
# directory = /code
# command=python manage.py sbch_blockheader_stream
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true

# [program:sbch__blocks_parser_queue]
# command = celery -A watchtower worker -n worker16 -l INFO -Ofair -Q sbch__blocks_parser_queue
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true

# [program:sbch__transactions_parser_queue]
# command = celery -A watchtower worker -n worker17 -l INFO -Ofair -Q sbch__transactions_parser_queue
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true

# [program:sbch__address_parser_queue]
# command = celery -A watchtower worker -n worker18 -l INFO -Ofair -Q sbch__address_parser_queue
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true

# [program:sbch__transaction_transfer_notification_queue]
# command = celery -A watchtower worker -n worker19 -l INFO -Ofair -Q sbch__transaction_transfer_notification_queue
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup=true